# target/linux/rockchip-rk3399/base-files/etc/init.d/fa-rk3399-netdev-cpu
#!/bin/sh /etc/rc.common

START=99
. /lib/functions.sh
start() {
  case "$(board_name)" in
  "friendlyelec,nanopi-r4s")
    # CPU1 (2nd core) -> internal dev
    INQETH0=$(awk -F":" "/eth0/ {print \$1}" </proc/interrupts | sed 's/\ //g')
    echo 2 > /proc/irq/${INQETH0}/smp_affinity
    # CPU2 (3nd core) -> pci-e dev
    INQETH1=$(awk -F":" "/eth1/ {print \$1}" </proc/interrupts | sed 's/\ //g')
    echo 4 > /proc/irq/${INQETH1}/smp_affinity
    ;;
  *)
    ;;
  esac
}
